# home/dot_config/systemd/user/atuin.service.d/20-distro-arch.conf
[Service]
ExecStart=
ExecStart=/usr/bin/atuin daemon
