# home/dot_config/systemd/user/atuin.service.d/70-device-dell-g15-linux.conf
[Service]
ExecStart=
ExecStart=%D/homebrew/opt/atuin/bin/atuin daemon
