# home/dot_config/autostart/1password.desktop
[Desktop Entry]
Type=Application
Name=1Password
Exec=1password --silent
StartupWMClass=1Password
X-Chezmoi=true
