# home/dot_config/systemd/user/atuin-daemon.socket
[Unit]
Description=atuin daemon socket

[Socket]
ListenStream=%t/atuin/atuin.sock
SocketMode=0700

[Install]
WantedBy=sockets.target
